# launchagents/com.local.postgresql14.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.local.postgresql14</string>
	<key>ProgramArguments</key>
	<array>
		<string>__REPO__/launchagents/run-service.sh</string>
		<string>postgresql14</string>
		<string>/opt/homebrew/opt/postgresql@14/bin/postgres</string>
		<string>-D</string>
		<string>/opt/homebrew/var/postgresql@14</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/postgresql14.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/postgresql14.error.log</string>
</dict>
</plist>
