# launchagents/com.local.cloudflared.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.local.cloudflared</string>
	<key>ProgramArguments</key>
	<array>
		<string>__REPO__/launchagents/run-service.sh</string>
		<string>cloudflared</string>
		<string>/opt/homebrew/bin/cloudflared</string>
		<string>tunnel</string>
		<string>run</string>
		<string>macbook</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/cloudflared.log</string>
	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/cloudflared.error.log</string>
</dict>
</plist>
